# build/supervisor/50-workers-forgejo.conf
# Workers are set to NOT auto start so we have time to enforce permissions on the cache first

[group:misp-workers]
programs=default,email,cache,prio,update,scheduler

[program:default]
directory=/workspace/MISP/MISP
command=/workspace/MISP/MISP/app/Console/cake start_worker default
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=false
autorestart=true
redirect_stderr=false
stderr_logfile=/tmp/misp-workers-errors.log
stdout_logfile=/tmp/misp-workers.log
directory=/workspace/MISP/MISP
user=www-data

[program:prio]
directory=/workspace/MISP/MISP
command=/workspace/MISP/MISP/app/Console/cake start_worker prio
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=false
autorestart=true
redirect_stderr=false
stderr_logfile=/tmp/misp-workers-errors.log
stdout_logfile=/tmp/misp-workers.log
directory=/workspace/MISP/MISP
user=www-data

[program:email]
directory=/workspace/MISP/MISP
command=/workspace/MISP/MISP/app/Console/cake start_worker email
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=false
autorestart=true
redirect_stderr=false
stderr_logfile=/tmp/misp-workers-errors.log
stdout_logfile=/tmp/misp-workers.log
directory=/workspace/MISP/MISP
user=www-data

[program:update]
directory=/workspace/MISP/MISP
command=/workspace/MISP/MISP/app/Console/cake start_worker update
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=false
autorestart=true
redirect_stderr=false
stderr_logfile=/tmp/misp-workers-errors.log
stdout_logfile=/tmp/misp-workers.log
directory=/workspace/MISP/MISP
user=www-data

[program:cache]
directory=/workspace/MISP/MISP
command=/workspace/MISP/MISP/app/Console/cake start_worker cache
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=false
autorestart=true
redirect_stderr=false
stderr_logfile=/tmp/misp-workers-errors.log
stdout_logfile=/tmp/misp-workers.log
user=www-data

[program:scheduler]
directory=/workspace/MISP/MISP
command=/workspace/MISP/MISP/app/Console/cake scheduler_worker
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=false
autorestart=true
redirect_stderr=false
stderr_logfile=/tmp/misp-workers-errors.log
stdout_logfile=/tmp/misp-workers.log
user=www-data